# docker/config/supervisord-headful.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:sway]
command=sway %(ENV_SWAY_EXTRA_ARGS)s
environment=
    XDG_RUNTIME_DIR="/tmp/xdg-runtime",
    WLR_BACKENDS="headless",
    WLR_LIBINPUT_NO_DEVICES="1",
    WLR_NO_HARDWARE_CURSORS="1",
    WLR_RENDERER="%(ENV_WLR_RENDERER)s",
    LIBSEAT_BACKEND="noop"
autostart=true
autorestart=true
stdout_logfile=/var/log/sway.log
stderr_logfile=/var/log/sway-err.log
priority=10

[program:wayvnc]
command=wayvnc --gpu --max-fps=60 --output=HEADLESS-1 0.0.0.0 %(ENV_VNC_PORT)s
environment=
    XDG_RUNTIME_DIR="/tmp/xdg-runtime",
    WAYLAND_DISPLAY="wayland-1"
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/wayvnc.log
stderr_logfile=/var/log/wayvnc-err.log
priority=20

[program:novnc]
command=websockify --web=/usr/share/novnc %(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc-err.log
priority=25

[program:chrome-1]
command=google-chrome-stable
    --no-sandbox
    --ozone-platform=wayland
    --enable-features=UseOzonePlatform
    --ignore-gpu-blocklist
    --enable-gpu-rasterization
    --enable-zero-copy
    --disable-dev-shm-usage
    --disable-background-networking
    --disable-component-update
    --disable-blink-features=AutomationControlled
    --disable-infobars
    --disable-session-crashed-bubble
    --no-first-run
    --no-default-browser-check
    --remote-debugging-port=%(ENV_CDP_PORT_1)s
    --user-data-dir=/tmp/chrome-profile-1
environment=
    XDG_RUNTIME_DIR="/tmp/xdg-runtime",
    WAYLAND_DISPLAY="wayland-1",
    DISPLAY=""
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/chrome-1.log
stderr_logfile=/var/log/chrome-1-err.log
priority=30

[program:chrome-2]
command=google-chrome-stable
    --no-sandbox
    --ozone-platform=wayland
    --enable-features=UseOzonePlatform
    --ignore-gpu-blocklist
    --enable-gpu-rasterization
    --enable-zero-copy
    --disable-dev-shm-usage
    --disable-background-networking
    --disable-component-update
    --disable-blink-features=AutomationControlled
    --disable-infobars
    --disable-session-crashed-bubble
    --no-first-run
    --no-default-browser-check
    --remote-debugging-port=%(ENV_CDP_PORT_2)s
    --user-data-dir=/tmp/chrome-profile-2
environment=
    XDG_RUNTIME_DIR="/tmp/xdg-runtime",
    WAYLAND_DISPLAY="wayland-1",
    DISPLAY=""
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/chrome-2.log
stderr_logfile=/var/log/chrome-2-err.log
priority=30
